# EasyTier/com.alick.swiftier.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alick.swiftier.helper</string>
    
    <!-- 可执行文件路径：使用 EasyTierHelper 作为产品名 -->
    <!-- 可执行文件路径：相对于 App Bundle 根目录 -->
    <key>BundleProgram</key>
    <string>Contents/MacOS/SwiftierHelper</string>
    
    <key>MachServices</key>
    <dict>
        <!-- Mach 服务名必须与代码中的 kHelperMachServiceName 一致 -->
        <key>com.alick.swiftier.helper</key>
        <true/>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/var/log/swiftier-helper-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/swiftier-helper-stderr.log</string>
    
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.alick.swiftier</string>
    </array>
</dict>
</plist>
